# infrastructure/launchd/com.digital-office.engine.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.digital-office.engine</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/God_Yurii/digital-office-engine/dist/index.js</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/God_Yurii/digital-office-engine</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/God_Yurii/digital-office-engine/shared/logs/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/God_Yurii/digital-office-engine/shared/logs/stderr.log</string>
</dict>
</plist>
